# launchd/com.allenneverland.cloudflared-tmux-chat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.allenneverland.cloudflared-tmux-chat</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{CLOUDFLARED_PATH}}</string>
        <string>tunnel</string>
        <string>run</string>
        <string>tmux-chat</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{LOG_DIR}}/cloudflared-tmux-chat.log</string>
    <key>StandardErrorPath</key>
    <string>{{LOG_DIR}}/cloudflared-tmux-chat.error.log</string>
</dict>
</plist>
